# scripts/systemd/vrising.service
[Unit]
Description=V Rising Dedicated Server (Baroboys)
Wants=network-online.target
After=network-online.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
Environment=DISPLAY=:0
WorkingDirectory=/home/bwinter_sc81/baroboys/VRising
User=bwinter_sc81
Group=bwinter_sc81

# Warn if already running, but always exit 0
ExecStartPre=/usr/bin/bash -c 'if pgrep -f VRisingServer.exe; then echo "🟡 Already running"; fi; exit 0'

# Use wrapper to detect intentional shutdown
ExecStart=/home/bwinter_sc81/baroboys/scripts/wrappers/vrising_launcher.sh

# No longer needed if wrapper handles exit logic, but harmless to keep
ExecStop=/usr/bin/pkill -f VRisingServer.exe

# Optional: avoid noisy crash logs if intentionally shutting down
SuccessExitStatus=0

StandardOutput=append:/home/bwinter_sc81/baroboys/VRising/logs/vrising.log
StandardError=append:/home/bwinter_sc81/baroboys/VRising/logs/vrising.log

KillMode=control-group
TimeoutStopSec=300
Restart=on-failure
