# scripts/services/xvfb/xvfb.service
[Unit]
Description=Xvfb headless X framebuffer for Wine
After=network.target

[Service]
ExecStart=/usr/bin/Xvfb :0 -screen 0 1024x768x16
User=bwinter_sc81
Group=bwinter_sc81
Restart=always

[Install]
WantedBy=multi-user.target
