# scripts/services/xvfb/xvfb-setup.service
[Unit]
Description=Xvfb Setup
Requires=network-online.target refresh-repo-startup.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root
User=root

ExecStart=/root/baroboys/scripts/services/xvfb/setup.sh

SuccessExitStatus=0
KillMode=control-group
TimeoutStartSec=300

StandardOutput=append:/var/log/baroboys/xfvb_setup.log
StandardError=append:/var/log/baroboys/xfvb_setup.log

[Install]
WantedBy=multi-user.target
