# scripts/services/xvfb/xvfb-startup.service
[Unit]
Description=Xvfb Startup
Requires=xvfb-setup.service
After=xvfb-setup.service

[Service]
Type=simple
Restart=always
WorkingDirectory=/home/bwinter_sc81
User=bwinter_sc81
Group=bwinter_sc81

ExecStart=/usr/bin/Xvfb :0 -screen 0 1024x768x24

# Type=simple means systemctl start returns as soon as Xvfb is forked — before it is
# actually listening on the display socket. ExecStartPost polls xdpyinfo until display
# :0 accepts connections, so any caller that does `systemctl start xvfb-startup.service`
# is guaranteed the display is ready when the command returns. Without this, wineboot
# races against a half-started Xvfb and fails to initialise the Wine prefix.
ExecStartPost=/bin/bash -c 'until xdpyinfo -display :0 &>/dev/null; do sleep 0.5; done'

SuccessExitStatus=0
KillMode=control-group
TimeoutStopSec=60

StandardOutput=append:/var/log/baroboys/xvfb_startup.log
StandardError=append:/var/log/baroboys/xvfb_startup.log

[Install]
WantedBy=multi-user.target
